# packaging/linux/deb/etc/xdg/autostart/eu-dss-agent.desktop
[Desktop Entry]
Type=Application
Name=EU-DSS Agent
Exec=/opt/eu-dss-agent/bin/eu-dss-agent
X-GNOME-Autostart-enabled=true
NoDisplay=true
